# docker/supervisord-api.conf
[program:bus-stats]
command=/api/bus-stats start
numprocs=1
autostart=true
autorestart=true